# mtda.service
[Unit]
Description=Multi-Tenant Device Access

[Service]
ExecStart=/usr/bin/mtda-cli -d

[Install]
WantedBy=multi-user.target
