# debian/mtda.postinst
#!/bin/sh

set -e

deb-systemd-helper enable mtda.service
grep -qe '^libcomposite$' /etc/modules || echo "libcomposite" >/etc/modules

