# debian/mtda.service
[Unit]
Description=Mentor Test Device Agent

[Service]
ExecStart=/usr/bin/mtda-cli -d

[Install]
WantedBy=multi-user.target
