# debian/mtda-systemd-helper.path
[Unit]
Wants=mtda.service

[Path]
PathChanged=/etc/mtda/config
PathChanged=/etc/mtda/config.d

[Install]
WantedBy=multi-user.target
