# debian/mtda-config.path
[Unit]
Wants=mtda.service

[Path]
PathChanged=/etc/mtda/config

[Install]
WantedBy=multi-user.target
